# linux-devkit/sysroots/cortexa8hf-vfp-neon-linux-gnueabi/var/lib/opkg/info/omap5-sgx-ddk-um-linux.list
/usr/share/sgx-lib/glsltest1_fragshaderB.txt
/usr/lib/gbm/gbm_pvr.so.1.9.6.0
/usr/lib/libusc.so.1
/usr/lib/libsrv_init.so.1
/etc/
/usr/lib/libPVRScopeServices.so.1.9.6.0
/usr/lib/libpvrDRMWSEGL.so.1.9.6.0
/usr/lib/libGLESv1_CM.so.1.9.6.0
/usr/lib/libsrv_init.so.1.9.6.0
/usr/lib/libsrv_um.so.1.9.6.0
/usr/lib/libpvr2d.so.1
/usr/bin/
/usr/lib/libpvr2d.so.1.9.6.0
/etc/init.d/pvr-init
/usr/lib/libpvrDRMWSEGL.so.1
/usr/lib/libpvrDRMWSEGL_FRONT.so.1.9.6.0
/etc/powervr.ini
/usr/lib/libglslcompiler.so.1.9.6.0
/usr/share/sgx-lib/glsltest1_vertshader.txt
/etc/init.d/
/usr/share/
/usr/share/sgx-lib/
/usr/lib/libsrv_um.so.1
/usr/lib/libglslcompiler.so.1
/usr/lib/libGLESv1_CM.so.1
/usr/lib/libEGL.so.1.9.6.0
/usr/lib/libusc.so.1.9.6.0
/usr/lib/libPVRScopeServices.so.1
/usr/lib/libIMGegl.so.1
/usr/lib/libGLES_CM.so.1
/usr/lib/libEGL.so.1
/usr/lib/gbm/
/usr/lib/gbm/gbm_pvr.so
/usr/share/sgx-lib/glsltest1_fragshaderA.txt
/usr/lib/libGLESv2.so.1.9.6.0
/usr/bin/pvrsrvinit
/usr/lib/libGLESv2.so.1
/usr/lib/libGLESv2.so.2
/usr/lib/libIMGegl.so.1.9.6.0
/usr/lib/libpvrDRMWSEGL_FRONT.so.1
